# .config/systemd/user/traefik.socket
[Unit]
Description=Traefik socket

[Socket]
ListenStream=192.168.3.1:80

[Install]
WantedBy=sockets.target